# launchd/com.frontierpapers.ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- 每晚 22:00 跑采集→排名→总结→推送微信（DeepSeek 耗时几分钟，推送约 22:02 到达）。
     安装：
       cp launchd/com.frontierpapers.ingest.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.frontierpapers.ingest.plist
     卸载：launchctl unload ~/Library/LaunchAgents/com.frontierpapers.ingest.plist
     注意：把下面的 WorkingDirectory 改成本机绝对路径。 -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.frontierpapers.ingest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "__PROJECT_DIR__" &amp;&amp; npm run ingest:send &gt;&gt; ingest.log 2&gt;&amp;1</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>22</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__PROJECT_DIR__/ingest.log</string>
  <key>StandardErrorPath</key>
  <string>__PROJECT_DIR__/ingest.log</string>
</dict>
</plist>
