# launchd/com.frontierpapers.refine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- 每周日 23:00 跑画像精炼，只生成建议 config/profile.suggested.md（不覆盖 profile.md）。
     安装：
       cp launchd/com.frontierpapers.refine.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.frontierpapers.refine.plist
     卸载：launchctl unload ~/Library/LaunchAgents/com.frontierpapers.refine.plist
     注意：下面路径若你的工程不在此处需相应修改。 -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.frontierpapers.refine</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "/Users/ali-jason/前沿论文情报台" &amp;&amp; npm run refine &gt;&gt; refine.log 2&gt;&amp;1</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>23</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/ali-jason/前沿论文情报台/refine.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ali-jason/前沿论文情报台/refine.log</string>
</dict>
</plist>
